# enterprise-attack/attack-pattern/attack-pattern--61afc315-860c-4364-825d-0d62b2e91edc.json
{
    "type": "bundle",
    "id": "bundle--f39ce5ae-ddfb-4e35-8482-f4a03ef27e26",
    "spec_version": "2.0",
    "objects": [
        {
            "x_mitre_platforms": [
                "Windows"
            ],
            "x_mitre_domains": [
                "enterprise-attack"
            ],
            "x_mitre_contributors": [
                "Scott Lundgren, @5twenty9, Carbon Black"
            ],
            "object_marking_refs": [
                "marking-definition--fa42a846-8d90-4e51-bc29-71d5b4802168"
            ],
            "type": "attack-pattern",
            "id": "attack-pattern--61afc315-860c-4364-825d-0d62b2e91edc",
            "created": "2020-01-24T15:51:52.317Z",
            "x_mitre_version": "1.0",
            "external_references": [
                {
                    "source_name": "mitre-attack",
                    "external_id": "T1547.003",
                    "url": "https://attack.mitre.org/techniques/T1547/003"
                },
                {
                    "source_name": "Github W32Time Oct 2017",
                    "url": "https://github.com/scottlundgren/w32time",
                    "description": "Lundgren, S. (2017, October 28). w32time. Retrieved March 26, 2018."
                },
                {
                    "source_name": "Microsoft W32Time May 2017",
                    "url": "https://docs.microsoft.com/windows-server/networking/windows-time-service/windows-time-service-tools-and-settings",
                    "description": "Mathers, B. (2017, May 31). Windows Time Service Tools and Settings. Retrieved March 26, 2018."
                },
                {
                    "source_name": "Microsoft W32Time Feb 2018",
                    "url": "https://docs.microsoft.com/windows-server/networking/windows-time-service/windows-time-service-top",
                    "description": "Microsoft. (2018, February 1). Windows Time Service (W32Time). Retrieved March 26, 2018."
                },
                {
                    "source_name": "Microsoft TimeProvider",
                    "url": "https://msdn.microsoft.com/library/windows/desktop/ms725475.aspx",
                    "description": "Microsoft. (n.d.). Time Provider. Retrieved March 26, 2018."
                },
                {
                    "source_name": "TechNet Autoruns",
                    "url": "https://technet.microsoft.com/en-us/sysinternals/bb963902",
                    "description": "Russinovich, M. (2016, January 4). Autoruns for Windows v13.51. Retrieved June 6, 2016."
                }
            ],
            "x_mitre_deprecated": false,
            "revoked": false,
            "description": "Adversaries may abuse time providers to execute DLLs when the system boots. The Windows Time service (W32Time) enables time synchronization across and within domains.(Citation: Microsoft W32Time Feb 2018) W32Time time providers are responsible for retrieving time stamps from hardware/network resources and outputting these values to other network clients.(Citation: Microsoft TimeProvider)\n\nTime providers are implemented as dynamic-link libraries (DLLs) that are registered in the subkeys of  <code>HKEY_LOCAL_MACHINE\\System\\CurrentControlSet\\Services\\W32Time\\TimeProviders\\</code>.(Citation: Microsoft TimeProvider) The time provider manager, directed by the service control manager, loads and starts time providers listed and enabled under this key at system startup and/or whenever parameters are changed.(Citation: Microsoft TimeProvider)\n\nAdversaries may abuse this architecture to establish persistence, specifically by registering and enabling a malicious DLL as a time provider. Administrator privileges are required for time provider registration, though execution will run in context of the Local Service account.(Citation: Github W32Time Oct 2017)",
            "modified": "2022-04-20T16:31:16.715Z",
            "created_by_ref": "identity--c78cb6e5-0c4b-4611-8297-d1b8b55e40b5",
            "name": "Time Providers",
            "x_mitre_detection": "Baseline values and monitor/analyze activity related to modifying W32Time information in the Registry, including application programming interface (API) calls such as <code>RegCreateKeyEx</code> and <code>RegSetValueEx</code> as well as execution of the W32tm.exe utility.(Citation: Microsoft W32Time May 2017) There is no restriction on the number of custom time providers registrations, though each may require a DLL payload written to disk.(Citation: Github W32Time Oct 2017)\n\nThe Sysinternals Autoruns tool may also be used to analyze auto-starting locations, including DLLs listed as time providers.(Citation: TechNet Autoruns)",
            "kill_chain_phases": [
                {
                    "kill_chain_name": "mitre-attack",
                    "phase_name": "persistence"
                },
                {
                    "kill_chain_name": "mitre-attack",
                    "phase_name": "privilege-escalation"
                }
            ],
            "x_mitre_is_subtechnique": true,
            "x_mitre_data_sources": [
                "Windows Registry: Windows Registry Key Modification",
                "Process: Process Creation",
                "Module: Module Load",
                "Command: Command Execution"
            ],
            "x_mitre_permissions_required": [
                "SYSTEM",
                "Administrator"
            ],
            "x_mitre_attack_spec_version": "2.1.0",
            "x_mitre_modified_by_ref": "identity--c78cb6e5-0c4b-4611-8297-d1b8b55e40b5"
        }
    ]
}