# enterprise-attack/attack-pattern/attack-pattern--a542bac9-7bc1-4da7-9a09-96f69e23cc21.json
{
    "type": "bundle",
    "id": "bundle--82ea37d4-48d0-4025-9803-2d5c7a8cccf2",
    "spec_version": "2.0",
    "objects": [
        {
            "type": "attack-pattern",
            "id": "attack-pattern--a542bac9-7bc1-4da7-9a09-96f69e23cc21",
            "created": "2020-10-12T17:50:31.584Z",
            "created_by_ref": "identity--c78cb6e5-0c4b-4611-8297-d1b8b55e40b5",
            "revoked": false,
            "external_references": [
                {
                    "source_name": "mitre-attack",
                    "url": "https://attack.mitre.org/techniques/T1053/006",
                    "external_id": "T1053.006"
                },
                {
                    "source_name": "Systemd Remote Control",
                    "description": "Aaron Kili. (2018, January 16). How to Control Systemd Services on Remote Linux Server. Retrieved July 26, 2021.",
                    "url": "https://www.tecmint.com/control-systemd-services-on-remote-linux-server/"
                },
                {
                    "source_name": "archlinux Systemd Timers Aug 2020",
                    "description": "archlinux. (2020, August 11). systemd/Timers. Retrieved October 12, 2020.",
                    "url": "https://wiki.archlinux.org/index.php/Systemd/Timers"
                },
                {
                    "source_name": "gist Arch package compromise 10JUL2018",
                    "description": "Catalin Cimpanu. (2018, July 10). ~x file downloaded in public Arch package compromise. Retrieved April 23, 2019.",
                    "url": "https://gist.github.com/campuscodi/74d0d2e35d8fd9499c76333ce027345a"
                },
                {
                    "source_name": "Arch Linux Package Systemd Compromise BleepingComputer 10JUL2018",
                    "description": "Catalin Cimpanu. (2018, July 10). Malware Found in Arch Linux AUR Package Repository. Retrieved April 23, 2019.",
                    "url": "https://www.bleepingcomputer.com/news/security/malware-found-in-arch-linux-aur-package-repository/"
                },
                {
                    "source_name": "acroread package compromised Arch Linux Mail 8JUL2018",
                    "description": "Eli Schwartz. (2018, June 8). acroread package compromised. Retrieved April 23, 2019.",
                    "url": "https://lists.archlinux.org/pipermail/aur-general/2018-July/034153.html"
                },
                {
                    "source_name": "Falcon Sandbox smp: 28553b3a9d",
                    "description": "Hybrid Analysis. (2018, July 11). HybridAnalsysis of sample 28553b3a9d2ad4361d33d29ac4bf771d008e0073cec01b5561c6348a608f8dd7. Retrieved September 8, 2023.",
                    "url": "https://www.hybrid-analysis.com/sample/28553b3a9d2ad4361d33d29ac4bf771d008e0073cec01b5561c6348a608f8dd7?environmentId=300"
                },
                {
                    "source_name": "Linux man-pages: systemd January 2014",
                    "description": "Linux man-pages. (2014, January). systemd(1) - Linux manual page. Retrieved April 23, 2019.",
                    "url": "http://man7.org/linux/man-pages/man1/systemd.1.html"
                }
            ],
            "object_marking_refs": [
                "marking-definition--fa42a846-8d90-4e51-bc29-71d5b4802168"
            ],
            "modified": "2025-04-15T19:58:55.648Z",
            "name": "Systemd Timers",
            "description": "Adversaries may abuse systemd timers to perform task scheduling for initial or recurring execution of malicious code. Systemd timers are unit files with file extension <code>.timer</code> that control services. Timers can be set to run on a calendar event or after a time span relative to a starting point. They can be used as an alternative to [Cron](https://attack.mitre.org/techniques/T1053/003) in Linux environments.(Citation: archlinux Systemd Timers Aug 2020) Systemd timers may be activated remotely via the <code>systemctl</code> command line utility, which operates over [SSH](https://attack.mitre.org/techniques/T1021/004).(Citation: Systemd Remote Control)\n\nEach <code>.timer</code> file must have a corresponding <code>.service</code> file with the same name, e.g., <code>example.timer</code> and <code>example.service</code>. <code>.service</code> files are [Systemd Service](https://attack.mitre.org/techniques/T1543/002) unit files that are managed by the systemd system and service manager.(Citation: Linux man-pages: systemd January 2014) Privileged timers are written to <code>/etc/systemd/system/</code> and <code>/usr/lib/systemd/system</code> while user level are written to <code>~/.config/systemd/user/</code>.\n\nAn adversary may use systemd timers to execute malicious code at system startup or on a scheduled basis for persistence.(Citation: Arch Linux Package Systemd Compromise BleepingComputer 10JUL2018)(Citation: gist Arch package compromise 10JUL2018)(Citation: acroread package compromised Arch Linux Mail 8JUL2018) Timers installed using privileged paths may be used to maintain root level persistence. Adversaries may also install user level timers to achieve user level persistence.(Citation: Falcon Sandbox smp: 28553b3a9d)",
            "kill_chain_phases": [
                {
                    "kill_chain_name": "mitre-attack",
                    "phase_name": "execution"
                },
                {
                    "kill_chain_name": "mitre-attack",
                    "phase_name": "persistence"
                },
                {
                    "kill_chain_name": "mitre-attack",
                    "phase_name": "privilege-escalation"
                }
            ],
            "x_mitre_attack_spec_version": "3.2.0",
            "x_mitre_contributors": [
                "SarathKumar Rajendran, Trimble Inc"
            ],
            "x_mitre_deprecated": false,
            "x_mitre_detection": "Systemd timer unit files may be detected by auditing file creation and modification events within the <code>/etc/systemd/system</code>, <code>/usr/lib/systemd/system/</code>, and <code>~/.config/systemd/user/</code> directories, as well as associated symbolic links. Suspicious processes or scripts spawned in this manner will have a parent process of \u2018systemd\u2019, a parent process ID of 1, and will usually execute as the \u2018root\u2019 user.\n\nSuspicious systemd timers can also be identified by comparing results against a trusted system baseline. Malicious systemd timers may be detected by using the systemctl utility to examine system wide timers: <code>systemctl list-timers \u2013all</code>. Analyze the contents of corresponding <code>.service</code> files present on the file system and ensure that they refer to legitimate, expected executables.\n\nAudit the execution and command-line arguments of the 'systemd-run' utility as it may be used to create timers.(Citation: archlinux Systemd Timers Aug 2020)",
            "x_mitre_domains": [
                "enterprise-attack"
            ],
            "x_mitre_is_subtechnique": true,
            "x_mitre_modified_by_ref": "identity--c78cb6e5-0c4b-4611-8297-d1b8b55e40b5",
            "x_mitre_platforms": [
                "Linux"
            ],
            "x_mitre_version": "1.3",
            "x_mitre_data_sources": [
                "File: File Modification",
                "Process: Process Creation",
                "Command: Command Execution",
                "Scheduled Job: Scheduled Job Creation"
            ]
        }
    ]
}