# enterprise-attack/attack-pattern/attack-pattern--f3c544dc-673c-4ef3-accb-53229f1ae077.json
{
    "type": "bundle",
    "id": "bundle--69dd8b52-1bf3-4536-94f0-8375b15d050c",
    "spec_version": "2.0",
    "objects": [
        {
            "x_mitre_data_sources": [
                "Process monitoring",
                "Process command-line parameters",
                "API monitoring"
            ],
            "x_mitre_permissions_required": [
                "User"
            ],
            "name": "System Time Discovery",
            "description": "The system time is set and stored by the Windows Time Service within a domain to maintain time synchronization between systems and services in an enterprise network. (Citation: MSDN System Time) (Citation: Technet Windows Time Service)\n\nAn adversary may gather the system time and/or time zone from a local or remote system. This information may be gathered in a number of ways, such as with [Net](https://attack.mitre.org/software/S0039) on Windows by performing <code>net time \\\\hostname</code> to gather the system time on a remote system. The victim's time zone may also be inferred from the current system time or gathered by using <code>w32tm /tz</code>. (Citation: Technet Windows Time Service) The information could be useful for performing other techniques, such as executing a file with a [Scheduled Task](https://attack.mitre.org/techniques/T1053) (Citation: RSA EU12 They're Inside), or to discover locality information based on time zone to assist in victim targeting.",
            "id": "attack-pattern--f3c544dc-673c-4ef3-accb-53229f1ae077",
            "x_mitre_platforms": [
                "Windows"
            ],
            "object_marking_refs": [
                "marking-definition--fa42a846-8d90-4e51-bc29-71d5b4802168"
            ],
            "x_mitre_version": "1.0",
            "type": "attack-pattern",
            "x_mitre_detection": "Command-line interface monitoring may be useful to detect instances of net.exe or other command-line utilities being used to gather system time or time zone. Methods of detecting API use for gathering this information are likely less useful due to how often they may be used by legitimate software.",
            "created_by_ref": "identity--c78cb6e5-0c4b-4611-8297-d1b8b55e40b5",
            "created": "2017-05-31T21:31:37.450Z",
            "kill_chain_phases": [
                {
                    "kill_chain_name": "mitre-attack",
                    "phase_name": "discovery"
                }
            ],
            "external_references": [
                {
                    "external_id": "T1124",
                    "source_name": "mitre-attack",
                    "url": "https://attack.mitre.org/techniques/T1124"
                },
                {
                    "source_name": "capec",
                    "external_id": "CAPEC-295",
                    "url": "https://capec.mitre.org/data/definitions/295.html"
                },
                {
                    "source_name": "MSDN System Time",
                    "description": "Microsoft. (n.d.). System Time. Retrieved November 25, 2016.",
                    "url": "https://msdn.microsoft.com/ms724961.aspx"
                },
                {
                    "source_name": "Technet Windows Time Service",
                    "description": "Mathers, B. (2016, September 30). Windows Time Service Tools and Settings. Retrieved November 25, 2016.",
                    "url": "https://technet.microsoft.com/windows-server-docs/identity/ad-ds/get-started/windows-time-service/windows-time-service-tools-and-settings"
                },
                {
                    "source_name": "RSA EU12 They're Inside",
                    "description": "Rivner, U., Schwartz, E. (2012). They\u2019re Inside\u2026 Now What?. Retrieved November 25, 2016.",
                    "url": "https://www.rsaconference.com/writable/presentations/file_upload/ht-209_rivner_schwartz.pdf"
                }
            ],
            "modified": "2019-06-18T14:07:20.239Z"
        }
    ]
}